# lib/crontab
SHELL=/bin/sh
PATH=/bin:/usr/bin:/usr/local/bin

30	*/4	*	*	*	-ns /usr/home/bin/fetchnews
